# .config/systemd/user/liminal-server.service
[Unit]
Description=Liminal Creative Space Server
After=network.target

[Service]
Type=simple
ExecStart=/home/liam/liminal/server/target/release/liminal-server
WorkingDirectory=/home/liam/liminal/server
Restart=always
RestartSec=3
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/liam/liminal/server

[Install]
WantedBy=default.target
